# supervisor.conf
[supervisord]
nodaemon=true

[program:flume]
command=/etc/init.d/flume-ng start
priority=10
